# Plugins/Apps/Tunnel/tunnelUDP.socket

[Unit]
Description=UDP tunnel listening on local port $PORT

[Socket]
ListenDatagram=$PORT

[Install]
WantedBy=sockets.target
