# Plugins/Apps/LetsEncrypt/renew-certificates.service

[Unit]
Description=Renew certificates

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot renew
